# .config/systemd/user/nm-applet.service
[Unit]
Description=Nm-Applet for persistence

[Service]
ExecStart=/usr/bin/nm-applet
Restart=always

[Install]
WantedBy=default.target
